# nanoclaw.service
[Unit]
Description=NanoClaw - Personal Claude Assistant
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=nanoclaw
WorkingDirectory=/home/nanoclaw/nanoclaw
ExecStart=/usr/bin/node /home/nanoclaw/nanoclaw/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

Environment="NODE_ENV=production"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target
